# task-3/crontabs
#2b -> Script untuk menjalankan tiap menit
* * * * * /home/helsasp/minute_log.sh

#2c -> Script untuk menjalankan tiap jam
0 * * * * /home/helsasp/aggregate_minutes_to_hourly_log.sh

