# task-3/crontabs
#3b -> Script untuk menjalankan tiap menit
* * * * * /home/helsasp/minute_log.sh

#3c -> Script untuk menjalankan tiap jam
0 * * * * /home/helsasp/aggregate_minutes_to_hourly_log.sh

#3d -> Zip tiap jam
0 * * * * /home/helsasp/backup_metrics.sh
