# task-4/crontabs
#4d -> Script untuk menjalankan tiap hari
0 0 * * * /home/helsasp/isabel-d.sh
