# task-4/crontabs
#zipping setiap 1 jam
0 * * * * /home/liares/isabel.sh

#4d -> Script untuk menjalankan tiap hari
0 0 * * * /home/helsasp/isabel.sh
