# task-2/crontabs
#2c
0 * * * * /home/liares/yu_database.sh
